# docker/cron/crontab
* * * * * php /var/www/html/artisan schedule:run >> /dev/stdout
